# deploy/mqtt-broker.service
[Unit]
Description=MQTT Broker Management Service
After=network.target mosquitto.service winter-river-hotspot.service
Requires=mosquitto.service winter-river-hotspot.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ECE-26.1-Winter-River
Environment="PATH=/home/pi/ECE-26.1-Winter-River/.venv/bin"
ExecStart=/home/pi/ECE-26.1-Winter-River/.venv/bin/python -m broker.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mqtt-broker-mgmt

[Install]
WantedBy=multi-user.target

# Instructions:
# 1. Copy this file to /etc/systemd/system/mqtt-broker.service
#    (setup_pi.sh does this automatically via sed substitution)
# 2. Run: sudo systemctl daemon-reload
# 3. Run: sudo systemctl enable mqtt-broker
# 4. Run: sudo systemctl start mqtt-broker
# 5. Check status: sudo systemctl status mqtt-broker
